# deploy/labot-poller.service
[Unit]
Description=labot.lv VPS agent poller (drains the intake queue, runs edit jobs)
After=network-online.target
Wants=network-online.target

[Service]
User=labot
WorkingDirectory=/srv/labot/labot-runtime
# Node parses the .env itself (proven-working invocation); absolute paths so cwd never matters.
ExecStart=/usr/bin/node --env-file=/srv/labot/labot-runtime/.env /srv/labot/labot-runtime/src/poll.mjs
Restart=always
RestartSec=5
# light hardening
NoNewPrivileges=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
